# mail/production/config/supervisor/postfix.ini
[program:postfix]
command=/bin/sh /usr/bin/postfix.sh
priority=2