# mail/production/config/supervisor/syslog.ini
[program:rsyslog]
command=rsyslogd -n
priority=1